# jenkins.sv.conf
[program:jenkins]
command=service jenkins start
redirect_stderr=true